# 0.10/root/etc/services.d/influxdb/run
#!/usr/bin/with-contenv sh

exec influxd -config=/usr/share/influxdb/config.toml
